# config/systemd/cv-backup.timer
[Unit]
Description=Catgirl Village Backup timer

[Timer]
OnCalendar=*-*-* 23:50:00
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
